# yfinance.service
[Unit]
Description=Yahoo Finance Python Microservice
Documentation=https://github.com/aristath/arduino-trader
After=network.target
Wants=network.target

[Service]
Type=simple
User=arduino
Group=arduino
WorkingDirectory=/home/arduino/app/repo/microservices/yfinance
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="PORT=9003"
EnvironmentFile=-/home/arduino/app/.env
ExecStart=/home/arduino/app/repo/microservices/yfinance/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 9003
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=yfinance

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

# Resource limits (lighter than trader service)
MemoryMax=512M
CPUQuota=25%

[Install]
WantedBy=multi-user.target
